# startup_sound.service
[Unit]
Description=startup sound

[Service]
ExecStart=/usr/bin/startup_sound

[Install]
WantedBy=multi-user.target

